# templates/boxtools/etc/network/if-up.d/reconfig2
#!/usr/bin/bash
echo "Waiting config script to finish!"
sleep 8
nohup /bin/bash /home/media/required.sh &
rm -- "$0"
exit 0
